# platform/linux/currently-fronting.service
[Unit]
Description=Currently Fronting

[Service]
Type=oneshot


